# etc/supervisor/conf.d/tm-api.conf
[program:tm-api]
command=/opt/tm/run/_api.sh
stdout_logfile=/opt/tm/log/api.log
autostart=true
autorestart=true
user=anonymous
stopsignal=KILL
numprocs=1
